# server/deploy/com.emergent-company.diane-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.emergent-company.diane-serve</string>

	<key>ProgramArguments</key>
	<array>
		<string>__BINARY_PATH__</string>
		<string>serve</string>
		<string>--api-port</string>
		<string>8890</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>WatchPath</key>
	<array>
		<string>__BINARY_PATH__</string>
	</array>

	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/diane-serve.log</string>

	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/diane-serve.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>WorkingDirectory</key>
	<string>__HOME__</string>

	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
